# djangowebsitet.conf
[program:gunicorn]
command=dbms_template_path/bin/start_gunicorn.sh
user=dbms_template_user
process_name=%(program_name)s
numproc=1
autostart=1
autorestart=1
redirect_stderr=true